# network-check.cron
*/10 * * * * root /home/pi/SugarPiDisplay/network-check.sh
